# udev/rules.d/99-framregs-dev.rules
ACTION=="add", SUBSYSTEM=="gpio", RUN+="/bin/sh -c 'chgrp -R gpio /sys/class/gpio/gpio*/ && chmod -R g+w /sys/class/gpio/gpio*/'"
